# system_files/usr/lib/systemd/system/bazzite-mx-groups.service
[Unit]
Description=bazzite-mx: append docker/libvirt groups and grant to wheel users
After=systemd-sysusers.service local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/bazzite-mx-groups
RemainAfterExit=yes

[Install]
WantedBy=default.target
